# crontab
SHELL=/bin/bash

# run every day at 03:40
40 3 * * * $HOME/openbylaws-saflii/fetch-all.sh >> $HOME/openbylaws-saflii/fetch-all.log 2>&1
