# crontab
SHELL=/bin/bash

# run every day at 03:40
40 3 * * * $HOME/saflii-acts-download/fetch-all.sh >> $HOME/saflii-acts-download/fetch-all.log 2>&1
